# scripts/certificates/com.certbot.renew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.certbot.renew</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/certbot</string>
        <string>renew</string>
        <string>--quiet</string>
        <string>--post-hook</string>
        <string>echo "Certbot renewal completed at $(date)" >> /var/log/certbot-renew.log</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Run twice daily at 00:00 and 12:00 -->
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/var/log/certbot-renew.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/certbot-renew.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
